# app/services/queue_executor_service/crontab
# START CRON JOB
* * * * * sleep 3; /usr/local/bin/python3 /src/app/services/queue_executor_service/main.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB